# photobooth-frontend.service
[Unit]
Description=Photobooth Frontend Kiosk
After=graphical-session.target
Wants=graphical-session.target
After=photobooth-backend.service

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
WorkingDirectory=/home/pi/photobooth
ExecStartPre=/bin/sleep 15
ExecStart=/home/pi/photobooth/start-kiosk.sh
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
