# photobooth-backend.service
[Unit]
Description=Photobooth Backend Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/photobooth/backend
ExecStart=/usr/bin/node server-raspberry.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

# GPIO Berechtigungen
SupplementaryGroups=gpio
DeviceAllow=/dev/gpiochip0 rw

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
